# x11-apps/radeon-profile-daemon/files/radeon-profile-daemon.initd
#!/sbin/openrc-run

description="Radeon Profile Daemon"

command="/usr/bin/radeon-profile-daemon"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need localmount
	use net
}
